# Testing/src/supervisor.conf
[program:myphp]
command=php /var/www/html/service.php